# mpc_framework/service/flask.conf
description "flask"
start on stopped rc RUNLEVEL=[2345]
respawn
exec python /home/koller/PycharmProjects/secretsharing/mpc_framework/run.py